# plexpy.conf
[program:plexpy]
user = nobody
command = /usr/bin/python /opt/Tautulli/Tautulli.py --nolaunch --datadir /config
