# plexpy.conf
[program:plexpy]
user = nobody
command = /usr/bin/python /opt/plexpy/PlexPy.py --nolaunch --datadir /config
